# systemd/localclaw-caddy.service
[Unit]
Description=LocalClaw Caddy reverse proxy
After=network-online.target localclaw-onebox.service
Wants=network-online.target localclaw-onebox.service

[Service]
Type=simple
User=spark
Group=spark
WorkingDirectory=/home/spark/src/dgx-spark-demo
Environment=HOME=/home/spark
ExecStart=/bin/bash /home/spark/src/dgx-spark-demo/scripts/start_dgx_caddy.sh
Restart=always
RestartSec=5
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
